# scripts/systemd/rahoot-deploy.timer
[Unit]
Description=Run the rahoot auto-deploy poller every 2 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=2min
Unit=rahoot-deploy.service

[Install]
WantedBy=timers.target
