# scripts/systemd/rahoot-deploy.service
[Unit]
Description=Rahoot auto-deploy poller (pull-based GitOps — deploy on new origin/main)
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/bash /nvmetank1/projects/rahoot/source/scripts/auto-deploy-poll.sh
# deploy.sh builds, smoke-tests, health-gates and auto-rolls-back — give it room
TimeoutStartSec=900
