# scripts/systemd/razzoozle-deploy.timer
[Unit]
Description=Run the razzoozle auto-deploy poller every 2 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=2min
Unit=razzoozle-deploy.service

[Install]
WantedBy=timers.target
