# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

# W5 rust-only cutover: the Node game socket is stopped — Rust (razzoozle-rust:3020)
# now serves socket + REST + result pages. nginx stays as the static SPA host and
# reverse-proxies all dynamic paths to Rust (docker/nginx.conf). The container's
# healthcheck (GET :3000/healthz) is nginx-proxied to Rust, so it stays green with
# no Node process. To restore the Node twin, re-add the [program:socket] block below.

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
